# systemd/jashect.service.d/local.conf
[Service]
Environment="USER=user"
Environment="PASS=pass"
